# wondershaper.service.in
[Unit]
Description=Bandwidth shaper/Network rate limiter
After=network-online.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/wondershaper -p
ExecStop=@sbindir@/wondershaper -p -c

[Install]
WantedBy=multi-user.target
